# quadlet/gitea-tailscale.container
[Unit]
Description=Tailscale sidecar for Matcha
After=network-online.target

[Container]
Image=tailscale/tailscale:latest
ContainerName=matcha-tailscale
Pod=matcha.pod

EnvironmentFile=__REPO__/.env
Environment=TS_STATE_DIR=/var/lib/tailscale
Environment=TS_SERVE_CONFIG=/config/tailscale-serve.json
Environment=TS_USERSPACE=false

Volume=%h/.local/share/matcha/tailscale:/var/lib/tailscale:Z
Volume=__REPO__/tailscale/config:/config
AddDevice=/dev/net/tun

AddCapability=NET_ADMIN
AddCapability=SYS_MODULE
SecurityLabelDisable=true

HealthCmd=tailscale status
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
