# quadlet/gitea-server.container
[Unit]
Description=Matcha gitea server
After=gitea-tailscale.service

[Container]
Image=gitea/gitea:1.25
ContainerName=matcha-gitea
Pod=matcha.pod

EnvironmentFile=__REPO__/.env
Environment=GITEA__security__INSTALL_LOCK=true
Environment=GITEA__database__DB_TYPE=sqlite3

Volume=%h/.local/share/matcha/gitea:/data:Z,U

HealthCmd=curl -fSs http://localhost:3000/api/healthz
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always

[Install]
WantedBy=default.target
